# src/etc/cron.d/alternc-acme
# Twice a day, every 12 hours :
# - renew the certificate of the panel (alternc-cerbot)
# - generate or renew the certificates of every domain available (generate_acme.php via alternc-acme)
0 */12 * * *  root /usr/lib/alternc/install.d/alternc-acme apache2 --quiet --all
